# Deployment_Doc/systemd/mqtt-timeline-alias.service
[Unit]
Description=MQTT Timeline mDNS Alias
After=network-online.target avahi-daemon.service
Wants=network-online.target avahi-daemon.service

[Service]
Type=oneshot
ExecStart=/usr/local/lib/mishka/mqtt-timeline-alias.sh
StandardOutput=journal
StandardError=journal
RemainAfterExit=yes
TimeoutStartSec=300
Restart=no

[Install]
WantedBy=multi-user.target
